# src/platform/linux/build-watcher.service
[Unit]
Description=Build Watcher MCP Server
After=network.target

[Service]
Type=simple
ExecStart=@@BINARY_PATH@@
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=build_watcher=info

[Install]
WantedBy=default.target
